# scripts/hil/debug/crontab
# wifi audit + tailscale (NO network_init restart — kills wpa every minute)
* * * * * sh /media/mmc/wifi_audit.sh
* * * * * sh /media/mmc/tailscale_wrapper.sh >> /media/mmc/tools.log 2>&1
